# packaging/sonos-sound-hub.service
[Unit]
Description=Sonos Sound Hub
Documentation=https://github.com/danmcpherson/SonosSoundHub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sonos-sound-hub
ExecStart=/opt/sonos-sound-hub/api
Restart=always
RestartSec=10
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=PATH=/home/danmc/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
